# config/com.glovenone3.ai-education-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.glovenone3.ai-education-weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/c8/workspace/ai-education-intelligence/scripts/weekly_codex_tick.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/c8/workspace/ai-education-intelligence/logs/launchd-weekly.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/c8/workspace/ai-education-intelligence/logs/launchd-weekly.err.log</string>
</dict>
</plist>
